# packages/kmod-6.18-nvidia-r580/load-grid-kernel-modules.service.in
[Unit]
Description=Load GRID kernel modules
RequiresMountsFor=PREFIX/lib/modules PREFIX/src/kernels
After=copy-grid-kernel-modules.service
Requires=copy-grid-kernel-modules.service
# Disable manual restarts to prevent loading kernel modules
# that weren't linked by the running system
RefuseManualStart=true
RefuseManualStop=true

[Service]
Type=oneshot
ExecCondition=/usr/bin/ghostdog match-nvidia-driver grid
ExecStart=/usr/bin/driverdog --modules-set nvidia-grid load-modules
RemainAfterExit=true
StandardError=journal+console

[Install]
RequiredBy=drivers.target
